# templates/default/profile.d/phpenv.sh.erb
export PHPENV_ROOT="/opt/phpenv"
export PATH="$PHPENV_ROOT/bin:$PATH"
eval "$(phpenv init -)"